# upstart/hashacat.conf
description "The hashacat.com flask server"

start on startup
stop on shutdown

script
	cd /home/john/cats
	exec /home/john/cats/server.py 2>&1 >> /var/log/hashacat.log
end script
